# files/turris_configured/etc/crontabs/root
MAILTO="petr.ruzicka@gmail.com"

*/5  * * * * vnstat -u
0    2 * * * /etc/graphs-vnstat.sh
