# files/openwrt/etc/crontabs/root
*/5 * * * * vnstat -u
0 2 * * * /etc/graphs-vnstat.sh
