# files/openwrt/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

echo -e "Subject: Reboot `uci get system.@system[0].hostname`.`uci get dhcp.@dnsmasq[0].domain`\n\nOpenwrt rebooted: `date; uptime`\n\n`grep -B 50 \"syslogd started\" /etc/messages`" | sendmail petr.ruzicka@gmail.com
exit 0
