# files/turris/etc/cron.d/nikola
MAILTO=""
*/15 * * * * root /usr/share/nikola/bin/nikola.sh
